# kbd-backlight/kbd-backlight.service
[Unit]
Description=Keyboard Backlight Service

[Service]
EnvironmentFile=-/etc/default/kbd-backlight
ExecStart=/usr/local/bin/kbd-backlight 2>&1

